# integration/services/skogai-lore-service.service
[Unit]
Description=SkogAI Lore Generator Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/skogix/skogai
ExecStart=/home/skogix/skogai/skogai-lore-service.sh llama3 300
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=OLLAMA_HOST=localhost:11434

[Install]
WantedBy=default.target